# deployment/s6-overlay/s6-rc.d/pulse-work/run
#!/usr/bin/with-contenv bash
set -e

cd ${ROOT}
echo "Starting Laravel Pulse worker..."
php artisan pulse:work