# deployment/s6-overlay/s6-rc.d/pulse-work/run
#!/command/execlineb
fdmove -c 2 1
with-contenv
cd /var/www/html

# Execute the Laravel Pulse worker process using pipeline and foreground for centralized logging
foreground {
    pipeline {
        php artisan pulse:work
    }
    /usr/local/bin/logger.sh 
    PULSE-WORK
}