# deployment/s6-overlay/s6-rc.d/pulse-work/run
#!/command/execlineb -P
fdmove -c 2 1
with-contenv
cd /var/www/html

# Log the command
/usr/local/bin/logger.sh PULSE-WORK "Starting Laravel Pulse worker"

# Execute the command and pipe all output to logger
pipeline { php artisan pulse:work }
/usr/local/bin/logger.sh PULSE-WORK