# deployment/s6-overlay/s6-rc.d/pulse-work/run
#!/command/execlineb -P
fdmove -c 2 1
s6-setuidgid nobody
cd ${ROOT}
php artisan pulse:work