# deployment/s6-overlay/etc/services.d/reverb/run
#!/usr/bin/with-contenv sh
exec php $ENV_ROOT/artisan reverb:start --host=127.0.0.1 --port=9000