# deployment/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb
fdmove -c 2 1
with-contenv
cd /var/www/html

# Execute the command and pipe all output to logger
pipeline { nginx } /usr/local/bin/logger.sh NGINX
