# deployment/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
fdmove -c 2 1
s6-setuidgid nobody
nginx -g "daemon off;"
