# deployment/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
fdmove -c 2 1
with-contenv
pipeline { nginx -g "daemon off;" }
/usr/local/bin/logger.sh NGINX INFO
