# deployment/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bash
set -e

echo "Starting Nginx server..."
exec nginx -g "daemon off;"
