# deployment/s6-overlay/s6-rc.d/horizon/run
#!/command/execlineb -P
fdmove -c 2 1
with-contenv
cd /var/www/html
pipeline { php artisan horizon }
/usr/local/bin/logger.sh HORIZON INFO