# deployment/s6-overlay/s6-rc.d/horizon/run
#!/command/execlineb
fdmove -c 2 1
with-contenv
cd /var/www/html

# Log the command
foreground {
    pipeline {
        php artisan horizon
    } 
    /usr/local/bin/logger.sh 
    HORIZON
}