# deployment/s6-overlay/s6-rc.d/horizon/run
#!/command/execlineb -P
fdmove -c 2 1

cd ${ROOT}
exec php artisan horizon