# deployment/s6-overlay/s6-rc.d/octane/run
#!/command/execlineb
fdmove -c 2 1
with-contenv
cd /var/www/html

# Import environment variables needed for the command
importas -i OCTANE_WORKERS OCTANE_WORKERS
importas -i OCTANE_TASK_WORKERS OCTANE_TASK_WORKERS
importas -i OCTANE_MAX_REQUESTS OCTANE_MAX_REQUESTS

# Execute the Laravel Octane server with continuous logging using pipeline
foreground {
    pipeline { 
        php artisan octane:start --server=swoole --host=0.0.0.0 --port=8000 --workers=${OCTANE_WORKERS} --task-workers=${OCTANE_TASK_WORKERS} --max-requests=${OCTANE_MAX_REQUESTS} 
    } 
    /usr/local/bin/logger.sh 
    OCTANE
}

