# deployment/s6-overlay/s6-rc.d/octane/run
#!/usr/bin/with-contenv bash
set -e

cd ${ROOT}
echo "Starting Laravel Octane server..."
php artisan octane:start --server=swoole --host=0.0.0.0 --port=8000 --workers=4 --task-workers=6 --max-requests=300