# deployment/s6-overlay/s6-rc.d/clear-scheduler-cache/run
#!/command/execlineb -P
fdmove -c 2 1
with-contenv
cd /var/www/html
pipeline { php artisan schedule:clear-cache }
/usr/local/bin/logger.sh SCHEDULER-CACHE INFO