# deployment/s6-overlay/s6-rc.d/clear-scheduler-cache/run
#!/usr/bin/with-contenv bash
set -e

cd /var/www/html
echo "Clearing Laravel scheduler cache..."
php artisan schedule:clear-cache