# deployment/s6-overlay/s6-rc.d/schedule-run/run
#!/command/execlineb -P
fdmove -c 2 1

cd ${ROOT}
php artisan schedule:work