# deployment/s6-overlay/s6-rc.d/schedule-run/run
#!/command/execlineb
fdmove -c 2 1
with-contenv
cd /var/www/html

# Execute the Laravel scheduler work process using pipeline for centralized logging
foreground { pipeline { php artisan schedule:work } /usr/local/bin/logger.sh SCHEDULER }