# deployment/s6-overlay/s6-rc.d/schedule-run/run
#!/command/execlineb -P
exec php $ENV_ROOT/artisan schedule:work