# deployment/s6-overlay/s6-rc.d/schedule-run/run
#!/command/execlineb -P
fdmove -c 2 1
with-contenv
cd /var/www/html
pipeline { php artisan schedule:work }
/usr/local/bin/logger.sh SCHEDULER INFO