# deployment/s6-overlay/s6-rc.d/schedule-run/run
#!/command/execlineb -P
fdmove -c 2 1
with-contenv
cd /var/www/html

# Log the command
/usr/local/bin/logger.sh SCHEDULER "Starting Laravel scheduler"

# Execute the command and pipe all output to logger
pipeline { php artisan schedule:work }
/usr/local/bin/logger.sh SCHEDULER