# deployment/s6-overlay/s6-rc.d/pulse-check/run
#!/command/execlineb -P
fdmove -c 2 1
with-contenv
cd /var/www/html

# Log the actual command with expanded variables
expanded_cmd="php artisan pulse:check"
/usr/local/bin/logger.sh PULSE-CHECK "${expanded_cmd}"

# Execute the command
pipeline { $expanded_cmd } /usr/local/bin/logger.sh PULSE-CHECK