# deployment/s6-overlay/s6-rc.d/pulse-check/run
#!/usr/bin/with-contenv bash
set -e

cd ${ROOT}
echo "Running Laravel Pulse checks..."
php artisan pulse:check