# deployment/s6-overlay/s6-rc.d/pulse-check/run
#!/command/execlineb
fdmove -c 2 1
with-contenv
cd /var/www/html

# Execute the Laravel Pulse check service using pipeline and foreground for centralized logging
foreground {
    pipeline {
        php artisan pulse:check
    } 
    /usr/local/bin/logger.sh 
    PULSE-CHECK
}