# deployment/s6-overlay/s6-rc.d/pulse-check/run
#!/command/execlineb -P
fdmove -c 2 1

cd /var/www/html
php artisan pulse:check